# backend/config/supervisor.conf
[program:odoo]
command=/usr/bin/odoo -c /etc/odoo/odoo.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/odoo.err.log
stdout_logfile=/var/log/odoo.out.log